# systemd/anglerfish-dashboard.service
[Unit]
Description=Anglerfish AI dashboard — operator UI + WebSocket stream
Documentation=https://github.com/retr0verride/Anglerfish-AI
After=network-online.target anglerfish-bridge.service
Wants=network-online.target

[Service]
Type=exec
User=anglerfish
Group=anglerfish
EnvironmentFile=/etc/anglerfish/anglerfish.env
# Closes pre-deploy sweep TODO-2: the previous "uvicorn --factory"
# invocation passed no args to create_app(settings) and so failed at
# startup on every clean install since Stage 4. The "anglerfish
# dashboard serve" subcommand owns its own uvicorn instance, loads
# settings explicitly, surfaces ValidationError as a clean exit 2,
# and mirrors the "bridge serve" + "lure serve" shape one-for-one.
# Host + port come from settings.dashboard.host/.port which the env
# file populates via ANGLERFISH_DASHBOARD__HOST/PORT.
ExecStart=/opt/anglerfish/venv/bin/anglerfish dashboard serve
Restart=on-failure
RestartSec=3s
TimeoutStopSec=10s

# Hardening — match the bridge + lure unit primitives so the dashboard
# does not lag behind on defense-in-depth. ReadWritePaths must include
# /var/log/anglerfish: the dashboard writes audit events from auth,
# routes, overrides_publisher, and audit_tailer to
# /var/log/anglerfish/audit.jsonl (AuditConfig default), which under
# ProtectSystem=strict would otherwise raise OSError(EROFS) on every
# login attempt + operator action. CapabilityBoundingSet retains
# CAP_NET_BIND_SERVICE because operators may rebind the dashboard to
# a privileged port via ANGLERFISH_DASHBOARD__PORT (parallels the
# lure unit's reasoning for port 22).
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
AmbientCapabilities=CAP_NET_BIND_SERVICE
ReadWritePaths=/var/lib/anglerfish /var/log/anglerfish
PrivateDevices=yes

[Install]
WantedBy=multi-user.target
