# systemd/anglerfish-firstboot.service
[Unit]
Description=Anglerfish AI first-boot configuration wizard
Documentation=https://github.com/retr0verride/Anglerfish-AI
After=network-online.target
ConditionPathExists=!/etc/anglerfish/anglerfish.env
DefaultDependencies=no
Before=anglerfish-bridge.service anglerfish-dashboard.service multi-user.target

[Service]
Type=oneshot
RemainAfterExit=no
StandardInput=tty
StandardOutput=tty
StandardError=tty
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
ExecStart=/opt/anglerfish/venv/bin/anglerfish-wizard --env /etc/anglerfish/anglerfish.env

# Re-enable getty@tty1 once the wizard completes successfully so the
# operator can log in on the console for diagnostics. The chroot hook
# masks getty@tty1 at image-build time so the wizard owns the tty on
# first boot; we undo that mask here. Only runs when ExecStart returned
# 0 (Type=oneshot's success contract).
ExecStartPost=/bin/systemctl unmask getty@tty1.service
ExecStartPost=/bin/systemctl start getty@tty1.service

# Schedule the MaxMind GeoLite2 update once the wizard has written the
# license key into the env file. The fetch unit no-ops if the key is
# absent (operator skipped that prompt).
ExecStartPost=/bin/systemctl start anglerfish-geo-update.service

# Lock the audit log to append-only at the filesystem level. The wizard
# has just written its first record so the file exists. `chattr +a` only
# works on ext2/3/4, btrfs and xfs; on other filesystems the command
# fails harmlessly and the audit log remains append-only by convention.
# Once set, the attribute persists in the inode across reboots — this
# only needs to fire on first boot.
#
# Leading `-` makes the whole step non-fatal: an unsupported filesystem
# or chmod failure must not abort firstboot.
ExecStartPost=-/bin/sh -c 'install -d -m 0750 -o anglerfish -g anglerfish /var/log/anglerfish && touch /var/log/anglerfish/audit.jsonl && chown anglerfish:anglerfish /var/log/anglerfish/audit.jsonl && chmod 0640 /var/log/anglerfish/audit.jsonl && chattr +a /var/log/anglerfish/audit.jsonl 2>/dev/null || true'

[Install]
WantedBy=multi-user.target
